# .config/systemd/user/lolcate.timer
[Unit]
Description=lolcate sync timer

[Timer]
OnCalendar=12:00:00
Unit=lolcate.service

[Install]
WantedBy=default.target
