# .config/systemd/user/offlineimap_all.timer
[Unit]
Description=OfflineIMAP Quicksync timer

[Timer]
OnCalendar=*:00:00
Unit=offlineimap_all.service

[Install]
WantedBy=default.target
